# wsl/cloud-drive-manager-mcp.service
[Unit]
Description=CloudDriveManager MCP Service
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/yuyutian/cloud-drive
ExecStart=/home/yuyutian/cloud-drive/bin/cloud-drive-mcp.sh
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target