# wsl/cloud-drive-manager-mcp.service
[Unit]
Description=CloudDriveManager MCP Service (HTTP Mode)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/CloudDriveManager
Environment="CLOUD_ENV=prod"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/CloudDriveManager/bin/cloud-drive-mcp --http --port 29313
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
LogsDirectory=cloud-drive-manager
RuntimeDirectory=cloud-drive-manager

[Install]
WantedBy=default.target
